# configs/systemd/mnt-downloads.automount
[Unit]
Description=Automount Downloads Volume

[Automount]
Where=/mnt/downloads
TimeoutIdleSec=900

[Install]
WantedBy=multi-user.target