# configs/systemd/mnt-media.automount
[Unit]
Description=Automount Media Volume

[Automount]
Where=/mnt/media
TimeoutIdleSec=900

[Install]
WantedBy=multi-user.target